# supervisor-deluged.conf
[program:deluged]
command=deluged -d -c /etc/deluge -L info -l /var/log/deluge/deluged.log

[program:deluge-web]
command=deluge-web -c /etc/deluge 
